# datasploit/service/celery/run
#!/bin/bash

C_FORCE_ROOT=root
cd /datasploit/core
celery -A core worker -l info --concurrency 20       

